# cronjob/files/etc/crontab
0 0 * * * /usr/bin/node /src/cli.js cache-list true >> /var/log/cron.log 2>&1